# cmd/dashboard/supervisord.conf
[supervisord]
childlogdir = /tmp

[program:dataloader]
command = /dataloader/bin/dataloader --db-host=skylb-dashboard-db --db-user=skylb_dashboard --db-passwd=vpnallowed --db-name=skylb_dashboard
redirect_stderr = true
stdout_logfile = /tmp/dataloader.log

[program:pgweb]
command = /dataloader/bin/pgweb --bind=0.0.0.0 -d
